# texpand.service
[Unit]
Description=texpand - Lightweight Wayland text expander
After=systemd-udev-settle.service

[Service]
ExecStart=%h/go/bin/texpand
Restart=on-failure
RestartSec=5
# Ensure wl-copy/wl-paste can reach the Wayland compositor
Environment=XDG_RUNTIME_DIR=/run/user/%U

[Install]
WantedBy=default.target
