# texpand.service
[Unit]
Description=texpand - Lightweight Wayland text expander
After=graphical-session.target systemd-udev-settle.service

# Allow more restart attempts during boot while devices settle
StartLimitIntervalSec=60
StartLimitBurst=10

[Service]
# Guard: don't start unless /dev/uinput is writable
ExecStartPre=/bin/sh -c 'test -w /dev/uinput'
ExecStart=%h/go/bin/texpand
Restart=on-failure
RestartSec=3
# Ensure wl-copy/wl-paste can reach the Wayland compositor
Environment=XDG_RUNTIME_DIR=/run/user/%U

[Install]
WantedBy=graphical-session.target
