# texpand.service
[Unit]
Description=texpand - Lightweight Wayland text expander
After=graphical-session.target

[Service]
ExecStart=%h/go/bin/texpand
Restart=on-failure
RestartSec=3
# Ensure wl-copy/wl-paste can reach the Wayland compositor
Environment=XDG_RUNTIME_DIR=/run/user/%U

[Install]
WantedBy=graphical-session.target
